# CassandraSeed.service
[Unit]
Description=Cassandra Seed
After=etcd.service
After=docker.service


[Service]
TimeoutStartSec=0
ExecStartPre= /usr/bin/docker pull sahardge/cassandra
ExecStart= /usr/bin/docker run -d --name seed sahardge/cassandra bash /usr/sbin/start-seed.sh


